# etc/init.d/laptop-mode
#!/bin/sh
#
# chkconfig: - 20 90
# description: Starts and stops "laptop-mode" - tweaks system behavior
#              to extend battery life.
#
# config:  /etc/laptop-mode/laptop-mode.conf

test -f /usr/sbin/laptop_mode || exit 0

if [ -f /lib/lsb/init-functions ] ; then
  . /lib/lsb/init-functions
else
  log_success_msg()
  {
    echo $1
  }
  log_failure_msg()
  {
    echo $1
  }
fi
TEMPFILE=`mktemp`

# Enable laptop mode when the system is booted when running on battery.

case $1 in
  start)
    mkdir -p /var/run/laptop-mode-tools
    touch /var/run/laptop-mode-tools/enabled
    if (/usr/sbin/laptop_mode auto > $TEMPFILE) ; then
      log_success_msg "$( cat $TEMPFILE )"
    else
      log_failure_msg "$( cat $TEMPFILE )"
    fi
    ;;

  restart|reload|force-reload)
    # Full restart: first stop laptop mode completely (to restore default mount options etc.)
    mkdir -p /var/run/laptop-mode-tools
    rm -f /var/run/laptop-mode-tools/enabled
    if (/usr/sbin/laptop_mode stop > $TEMPFILE) ; then
      log_success_msg "$( cat $TEMPFILE )"
    else
      log_failure_msg "$( cat $TEMPFILE )"
    fi

    # Now remove files containing stored status, re-enable, and start it up again.
    rm -f /var/run/laptop-mode-tools/*
    touch /var/run/laptop-mode-tools/enabled
    if (/usr/sbin/laptop_mode auto force > $TEMPFILE) ; then
      log_success_msg "$( cat $TEMPFILE )"
    else
      log_failure_msg "$( cat $TEMPFILE )"
    fi
    ;;

  stop)
    rm -f /var/run/laptop-mode-tools/enabled
    if (/usr/sbin/laptop_mode stop > $TEMPFILE) ; then
      log_success_msg "$( cat $TEMPFILE )"
    else
      log_failure_msg "$( cat $TEMPFILE )"
    fi
    ;;

  status)
    echo "Laptop mode status:"
    echo
    /usr/sbin/laptop_mode status
    ;;
    
  *)
    echo "Usage: $0 {stop|start|restart|reload|force-reload|status}" >&2
    exit 2
    ;;
esac

if [ -f $TEMPFILE ] ; then
	rm $TEMPFILE
fi

exit 0
