# etc/init.d/laptop-mode
#!/bin/sh
#
# chkconfig: - 20 90
# description: Starts and stops "laptop-mode" - tweaks system behavior
#              to extend battery life.
#
# config:  /etc/laptop-mode/laptop-mode.conf

### BEGIN INIT INFO
# Provides:          laptop-mode
# Should-Start:      $all
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Enable laptop-mode-tools power management functions
# Description:       Enable laptop-mode-tools power management functions
### END INIT INFO


test -f /usr/sbin/laptop_mode || exit 0

. /lib/lsb/init-functions

# Enable laptop mode when the system is booted when running on battery.

TEMPFILE=`mktemp`
echo "Failed to start laptop mode" > $TEMPFILE


trap "log_action_end_msg 10 $( cat $TEMPFILE ) ; rm $TEMPFILE; exit 10" HUP INT ABRT QUIT SEGV TERM
set -e

case $1 in
  start)
    log_action_begin_msg "Enabling laptop mode"
    touch /var/run/laptop-mode-enabled
    /usr/sbin/laptop_mode auto > $TEMPFILE
    log_action_end_msg 0 $( cat $TEMPFILE )
    ;;

  restart|reload|force-reload)
    log_action_begin_msg "Disabling laptop mode"
    # Full restart: first stop laptop mode completely (to restore default mount options etc.)
    rm -f /var/run/laptop-mode-enabled
    /usr/sbin/laptop_mode stop > $TEMPFILE
    log_action_end_msg 0 $( cat $TEMPFILE )
    
    log_action_begin_msg "Enabling laptop mode"
    # Now remove files containing stored status, re-enable, and start it up again.
    rm -f /var/run/laptop-mode-*
    touch /var/run/laptop-mode-enabled
    /usr/sbin/laptop_mode auto force > $TEMPFILE
    log_action_end_msg 0 $( cat $TEMPFILE )
    ;;

  stop)
    log_action_begin_msg "Disabling laptop mode"
    rm -f /var/run/laptop-mode-enabled
    /usr/sbin/laptop_mode stop > $TEMPFILE
    log_action_end_msg 0 $( cat $TEMPFILE )    
    ;;

  status)
    echo "Laptop mode status:"
    echo
    /usr/sbin/laptop_mode status
    ;;
    
  *)
    echo "Usage: $0 {stop|start|restart|reload|force-reload|status}" >&2
    exit 2
    ;;
esac
exit 0
