# etc/rules/99-laptop-mode.rules
SUBSYSTEM=="usb", ACTION=="add", RUN+="/usr/sbin/laptop_mode force modules=usb-autosuspend devices=%k"
